# layers/meta-bbb/recipes-bbb/hmi/files/hmi.service
[Unit]
Description=Continuous Hello, HMI Service

[Service]
ExecStart=/bin/bash  /usr/bin/hmi.sh

[Install]
WantedBy=multi-user.target